# files/flatpak/flatpak-setup.timer
[Unit]
Description=Timer for flatpak-setup

[Timer]
OnBootSec=30

[Install]
WantedBy=timers.target
